# casper/debian/casper0.postinst
# Automatically added by dh_installinit
if [ -x "/etc/init.d/casper0" ] || [ -e "/etc/init/casper0.conf" ]; then
	if [ ! -e "/etc/init/casper0.conf" ]; then
		update-rc.d casper0 defaults >/dev/null
	fi
	invoke-rc.d casper0 start || exit $?
fi
# End automatically added section
# Automatically added by dh_makeshlibs
if [ "$1" = "configure" ]; then
	ldconfig
fi
# End automatically added section
